# supervisord.conf
[supervisord]
nodaemon=true

[program:influx]
command=/entrypoint.sh

[program:ra-stats]
command=/stats/gathering-script.sh >> /dev/stdout
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0