# supervisord.conf
[supervisord]
nodaemon=true

[program:influx]
command=/entrypoint.sh

[program:ra-stats]
command=/stats/script.sh >> /dev/stdout